# config/btc-collector.service
[Unit]
Description=Bitcoin Data Collector Service
After=network.target

[Service]
User=%i
WorkingDirectory=/home/%i/PycharmProjects/btc
ExecStart=/home/%i/PycharmProjects/btc/.venv/bin/python /home/%i/PycharmProjects/btc/scripts/btc_collector_service.py
Restart=on-failure
RestartSec=60s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btc-collector
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 